# macos/Library/LaunchAgents/se.victorandree.KeyMap.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
This launch agent uses the `hidutil` utility to remap keys.

- Grave Accent and Tilde (0x35) <-> Non-US \ and | (0x64)

See https://www.nanoant.com/mac/macos-function-key-remapping-with-hidutil
-->
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>se.victorandree.KeyMap</string>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/hidutil</string>
      <string>property</string>
      <string>--matching</string>
      <string>{"VendorID":1452,"ProductID":630}</string>
      <string>--set</string>
      <string>{"UserKeyMapping":[{"HIDKeyboardModifierMappingSrc":30064771172,"HIDKeyboardModifierMappingDst":30064771125},{"HIDKeyboardModifierMappingSrc":30064771125,"HIDKeyboardModifierMappingDst":30064771172}]}</string>
    </array>
    <key>RunAtLoad</key>
    <true />
  </dict>
</plist>
